# bluetooth_2_usb.service
[Unit]
Description=Bluetooth-2-USB HID relay
After=bluetooth.service network-online.target local-fs.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/bluetooth_2_usb
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/default/bluetooth_2_usb
ExecStart=/bin/sh -lc 'exec "/opt/bluetooth_2_usb/venv/bin/python" -m bluetooth_2_usb ${BLUETOOTH_2_USB_ARGS:---auto_discover --grab_devices --interrupt_shortcut CTRL+SHIFT+F12}'
Restart=on-failure
RestartSec=2s
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
StateDirectory=bluetooth_2_usb
LogsDirectory=bluetooth_2_usb
RuntimeDirectory=bluetooth_2_usb
ReadWritePaths=/var/lib/bluetooth

[Install]
WantedBy=multi-user.target
